# dz/roles/dbserver/files/db-backup.timer
[Unit]
Description=Borg backup timer

[Timer]
OnBootSec=1min
OnUnitActiveSec=5min
Unit=db-backup.service

[Install]
WantedBy=multi-user.target
